# images/portal/resources/supervisor/portal.conf
[program:portal]
command=python /opt/portal/pod/pod.py
