# images/scaler/resources/supervisor/scaler.conf
[program:scaler]
command=python /opt/scaler/pod/pod.py
